# cloud/fake-review-supervisor.conf
# Install as /etc/supervisor/conf.d/fake-review.conf on the Vast.ai instance,
# then run: supervisorctl reread && supervisorctl update
[program:fake-review]
directory=/workspace/Fake-review
command=/bin/bash /workspace/Fake-review/cloud/run_supervised.sh
environment=PROC_NAME="%(program_name)s",SCOPE="modern",LANES="1"
autostart=true
autorestart=false
startsecs=5
startretries=1
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
