# cloud/fake-review-publish-supervisor.conf
[program:fake-review-publish]
directory=/workspace/Fake-review
command=/workspace/Fake-review/cloud/publish_adapters.sh
environment=PROC_NAME="%(program_name)s",HF_REPO_ID="KandirResearch/fake-review-rs-qlora-adapters",HF_HOME="/workspace/.hf_home"
autostart=false
autorestart=unexpected
startsecs=5
startretries=1
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/Fake-review/results/publish_adapters.log
redirect_stderr=true
stdout_logfile_maxbytes=0
